# ezs.service
[Unit]
Description=eazy sessions session/display manager
After=graphical.target

[Service]
ExecStart=/usr/bin/ezs
StandardInput=tty
TTyPath=/dev/tty1
StandardOutput=tty

[Install]
Alias=display-manager.service
WantedBy=multi-user.target