# ezs.service
[Unit]
Description=ezs
After=getty@tty1.service

[Service]
ExecStart=/usr/bin/ezs-dm
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty1
TTYReset=yes
TTYVHangup=yes
Type=idle

[Install]
Alias=display-manager.service