# ezs.service
[Unit]
Description=ezs
After=multi-user.target
Requires=multi-user.target

[Service]
ExecStart=/usr/bin/ezs-dm
StandardInput=tty
StandardOutput=tty
TTYPath=/dev/tty1
Restart=always

[Install]
WantedBy=default.target
Alias=display-manager.service